# infra/proxmox/ansible/roles/github-runner/templates/runner-cleanup.timer.j2
[Unit]
Description=GitHub Actions Runner Cleanup (every 6 hours)
Documentation=https://github.com/benjr70/Smart-Smoker-V2

[Timer]
# Run every 6 hours (at 00:00, 06:00, 12:00, 18:00)
OnCalendar=*-*-* 00/6:00:00
RandomizedDelaySec=15min
# Run immediately if last scheduled run was missed
Persistent=true
# Run shortly after boot if timer was missed
OnBootSec=15min

[Install]
WantedBy=timers.target

